# alpine-php7/config/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/var/logs/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:php-fpm]
command=php-fpm7 -F
