# alpine-redis/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/var/logs/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D

[program:redis]
command=/usr/local/bin/redis-server /etc/redis/redis.conf

